# com.smart_money_tracker.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.smart_money_tracker.monitor</string>

  <key>WorkingDirectory</key>
  <string>/Users/wj/Documents/smart_money_tracker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/wj/Documents/smart_money_tracker/run_monitor.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/wj/Documents/smart_money_tracker/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/wj/Documents/smart_money_tracker/launchd.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
